# scaleout/slurmd.service
[Service]
ExecCondition=bash -c 'test $(hostname -s | cut -b1-4) = "node"'
ExecStartPre=/usr/local/bin/slurmd.startup.sh
